# deploy/phonecost-backend.service
[Unit]
Description=PhoneCost Backend
After=mysqld.service redis.service
Requires=mysqld.service redis.service

[Service]
Type=simple
User=admin
WorkingDirectory=/data/apps/phonecost/backend/phonecost
ExecStart=/usr/bin/java -jar target/phonecost-0.1.0-SNAPSHOT.jar --spring.profiles.active=prod
Restart=always
RestartSec=10
# 以下环境变量必填，安装时由脚本自动替换或手动填写
Environment=DB_HOST=localhost
Environment=DB_PORT=3306
Environment=DB_NAME=phonecost
Environment=DB_USER=phonecost_app
Environment=DB_PASSWORD=PLACEHOLDER_DB_PASSWORD
Environment=REDIS_HOST=localhost
Environment=REDIS_PORT=6379
Environment=REDIS_PASSWORD=PLACEHOLDER_REDIS_PASSWORD
Environment=JWT_SECRET=PLACEHOLDER_JWT_SECRET

[Install]
WantedBy=multi-user.target
